# deployments/ansible/roles/node-agent/templates/node-agent.service.j2
[Unit]
Description=HA Node Agent for {{ node_id }}
After=postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=postgres
Environment=NODE_ID={{ node_id }}
Environment=NODE_ADDR={{ node_addr }}
Environment=ORCHESTRATOR_ADDR={{ orchestrator_addr }}
Environment=PG_USER={{ pg_user }}
Environment=PG_PASSWORD={{ pg_password }}
Environment=PGDATA=/var/lib/postgresql/{{ pg_version }}/main
ExecStart=/usr/local/bin/node-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
